# systemd/dropins/gandalf-backup.service.d/20-clean-lock.conf
[Service]
ExecStopPost=/bin/rm -f /run/gandalf-backup.lock
