# systemd/dropins/gandalf-backup.service.d/override.conf
[Service]
SyslogIdentifier=gandalf-backup
StandardOutput=journal
StandardError=journal
