# systemd/dropins/gandalf-healthcheck.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/bash /usr/local/sbin/gandalf-healthcheck
